# ansible/webapp/roles/app/templates/supervisor-dotnet.conf
[program:webapi]
command=/usr/bin/dotnet /home/app/publish/MvcSample.dll
directory=/home/app/publish
autostart=true
autorestart=true